# systemd/aesop-ui.service
[Unit]
Description=Aesop's UI service
After=network.target

[Service]
ExecStart=/usr/bin/python3 -m aesop.ui

[Install]
WantedBy=multi-user.target
